# ogd-rabbitmq-worker.conf
###
# Installation:
# -------------
# Copy this file to /etc/init/
# sudo cp --force /var/www/ogd-rabbitmq-worker.conf /etc/init/ogd-rabbitmq-worker.conf
#
# Usage:
# ------
# Start deamon:
# sudo service ogd-rabbitmq-worker start
# Stop deamon:
# sudo service ogd-rabbitmq-worker stop
# Getting status of process:
# sudo service ogd-rabbitmq-worker status
###

# Info
description "OGD RabbitMQ Worker"
author      "Liip AG / Team Jazz / Juerg Hunziker"

# Events
start on startup
stop on shutdown

# Automatically respawn
respawn
respawn limit 20 5

# Run the script!
# Note, in this example, if your PHP script returns
# the string "ERROR", the daemon will stop itself.
script
    [ $(exec /usr/bin/php -f /var/www/ogdch.dev/rabbitmq/read.php) = 'ERROR' ] && ( stop; exit 1; )
end script

